# callweaver-1.2.1/contrib/init.d/.svn/text-base/rc.redhat.callweaver.svn-base
#!/bin/sh
# $Id: callweaver,v 0.2 2004/07/18 20:24:07 Gregory Boehnlein <damin@nacs.net>
#
# callweaver    Starts, Stops and Reloads CallWeaver.
#
# chkconfig: 2345 40 60
# description: CallWeaver and telephony daemon.
# processname: callweaver
# pidfile: /var/run/callweaver.pid
#
# Sun Jul 18 2004 Gregory Boehnlein <damin@nacs.net>
# - Updated Version to 0.2
# - Added test for safe_callweaver
# - Verified SIGTERM issued by "killproc" ensures "stop gracefully"
# - Added support for -U and -G command line options
# - Modified "reload" to call callweaver -rx 'reload' 

# Source function library.
. /etc/rc.d/init.d/functions

[ -x /usr/sbin/callweaver ] || exit 0
[ -d /etc/callweaver ] || exit 0

# Leave this set unless you know what you are doing.
export LD_ASSUME_KERNEL=2.4.1

# Full path to callweaver binary
DAEMON=/usr/sbin/callweaver

# Full path to safe_callweaver script
SAFE_CALLWEAVER=/usr/sbin/safe_callweaver

# Uncomment the following and set them to the user/groups that you
# want to run CallWeaver as. NOTE: this requires substantial work to
# be sure that CallWeaver's environment has permission to write the
# files required  for  its  operation, including logs, its comm
# socket, the callweaver database, etc.
#AST_USER="callweaver"
#AST_GROUP="callweaver"

RETVAL=0

start() {
	# Start daemons.
	echo -n $"Starting callweaver: "
        if [ -f $SAFE_CALLWEAVER ] ; then
		DAEMON=$SAFE_CALLWEAVER
	fi
	if [ $AST_USER ] ; then
		ASTARGS="-U $AST_USER"
	fi
	if [ $AST_GROUP ] ; then
		ASTARGS="`echo $ASTARGS` -G $AST_GROUP"
	fi
	daemon $DAEMON $ASTARGS
	RETVAL=$?
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/callweaver
	echo
	return $RETVAL
}

stop() {
	# Stop daemons.
	echo -n $"Shutting down callweaver: "
	killproc callweaver
	RETVAL=$?
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/callweaver
	echo
	return $RETVAL
}

restart() {
	stop
	start
}

reload() {
	$DAEMON -rx 'reload' > /dev/null 2> /dev/null
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart)
	restart
	;;
  reload)
	reload
	;;
  condrestart)
	[ -f /var/lock/subsys/callweaver ] && restart || :
	;;
  status)
	status callweaver
	;;
  *)
	echo "Usage: callweaver {start|stop|restart|reload|condrestart|status}"
	exit 1
esac

exit $?
